# debian/preinst
#! /bin/bash

#
# Preinst setups

update-rc.d  -f gosocksd remove
